# quadlet/git-friends-server.container
# Git Friends Server Quadlet

[Unit]
Description=Git Friends Server
Documentation=https://github.com/containers/podman/blob/main/docs/source/markdown/podman-systemd.unit.5.md
After=git-friends-network.service git-friends-mqtt.service
Requires=git-friends-network.service
Wants=git-friends-mqtt.service git-friends-logs.service

[Container]
ContainerName=git-friends-server
Image=localhost/git-friends:latest
PublishPort=8080:8080
Volume=%i/docker/config/git-friends.toml:/app/config/git-friends.toml:ro,Z
Volume=git-friends-logs:/app/logs:Z
Environment=RUST_LOG=info
Environment=GIT_FRIENDS_CONFIG=/app/config/git-friends.toml
Network=git-friends-network.network
Exec=./bin/gf-server --bind 0.0.0.0:8080
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

[Service]
Restart=unless-stopped
TimeoutStartSec=900

[Install]
WantedBy=default.target
